# alarming-protocol/files/alarming-protocol.init
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
	if[ "$(uci get seismo.block.alarming_protocol)" = "0"]
	then
		cd /data/
		/usr/bin/alarming-protocol $1 $2 >>/data/alarming-protocol.log 2>>/data/alarming-protocol.errlog &
	fi
}

stop() {
	if[ "$(uci get seismo.block.alarming_protocol)" = "0"]
	then
		killall alarming-protocol
		killall -9 alarming-protocol
	fi
}


