# node-info-service/files/node-info-service.init
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
	/usr/bin/node-info-service &
}

stop() {
	killall node-info-service
}

