# dlr-bluetooth/files/dlr-bluetooth.init
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
	cd /usr/bin/dlr-bluetooth
	# wipe old pids of scripts
	> .pid
	for script in *.sh; do
		logger -t DLR_BT "starting $PWD/$script"
		./$script &
		echo $! >> .pid
	done
}
   
stop() {
       pkill jamvm
       pkill hci*
}