# seismo-misc/files/startup_node_check.initd
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=99
start() {
	if[ "$(uci get seismo.block.startup_node_check)" = "0"]
	then
		/sbin/startup_node_check.sh > /dev/null 2>&1 < /dev/null &
	fi
}
