# wifiwatchdog/files/wifiwatchdog.initd
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org
# Somehow priority option in watchdog.conf does not work, so its started with nice and prio -1

START=99
start() {
	if[ "$(uci get seismo.block.wifiwatchdog)" = "0"]
	then
		/usr/sbin/wifiwatchdog.sh start > /dev/null 2>&1 < /dev/null &
	fi
}
